# install/all_in_one/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod

[program:nodeapp]
command=npm start --production

[program:admin]
command=npm run admin
